# supervisord.conf
[program:node]
command = node server.js
directory = /home/dotcloud/current

[program:newrelic]
command=/home/dotcloud/current/nrsysmond.x64 -f